# birdsong-tft.service
# systemd user service for the Mini PiTFT companion screen (tft.py).
#
# Install (on the Pi, after enabling SPI + spidev.bufsiz — see README):
#   cp birdsong-tft.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now birdsong-tft.service
#
# Requires birdsong.service (it reads localhost:8000) and linger enabled.

[Unit]
Description=Birdsong TFT companion display
# No After=birdsong.service — birdsong.service is After=default.target, which
# would form an ordering cycle and silently drop this unit at boot. tft.py
# tolerates the web app not being ready yet (it polls with timeouts).

[Service]
WorkingDirectory=%h/birdsong
ExecStart=%h/birdsong/.venv/bin/python -u %h/birdsong/tft.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
